# 8.3/test/supervisord-with-octane.conf
[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "while [ true ]; do (php /var/www/site/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/site/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/site/worker.log
stopwaitsecs=3600


[program:laravel-octane]
# roadrunner
# command=/usr/bin/php -d variables_order=EGPCS /var/www/site/artisan octane:start --server=roadrunner --host=0.0.0.0 --rpc-port=6001 --port=8000
# swoole
# command=/usr/bin/php -d variables_order=EGPCS /var/www/site/artisan octane:start --server=swoole --host=0.0.0.0 --port=8000
#
command=/usr/bin/php -d variables_order=EGPCS /var/www/site/artisan octane:start --host=0.0.0.0 --port=8000 --max-requests=250
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
